# scripts/com.lybrish.pdf-download.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lybrish.pdf-download</string>

  <!-- 直接呼叫 venv python，不經過 bash（同 KB-Scripts 模式）。
       原因：macOS TCC 追蹤 binary identity；bash 是系統共用的，難給單獨 FDA。
       venv python 是 symlink 到 /Library/.../python3，跟 KB 共用同一個 binary，
       借用 KB 已給的 FDA 授權。 -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/lybrish-mini/lybrish_claude/projects/pdf-download/venv/bin/python</string>
    <string>-m</string>
    <string>pdf_download.cli</string>
    <string>fetch</string>
    <string>--notify</string>
    <string>--silent-when-empty</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/lybrish-mini/lybrish_claude/projects/pdf-download</string>

  <!-- 排程：每週日 16:00 (Asia/Taipei) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>16</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/lybrish-mini/Library/Logs/pdf-download/fetch-auto-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/lybrish-mini/Library/Logs/pdf-download/fetch-auto-stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/lybrish-mini</string>
  </dict>

  <key>StartCalendarIntervalAllowsMultipleEvents</key>
  <false/>
</dict>
</plist>
